# src/server/etc/cron.d/hostdb
0 0 * * * www-data /usr/local/bin/hostdb_git_gc.sh
